# bbd_server/docker_assets/supervisor_confs/apache2.conf
[program:apache2]
command=/bin/bash -c "service apache2 start"
redirect_stderr=true
