# bbd_server/docker_assets/supervisor_confs/apache2.conf
[program:apache2]
command=apache2ctl -DFOREGROUND
redirect_stderr=true
autorestart=true
