# bbd_server/docker_assets/supervisor_confs/bbd_daemon.conf
[program:bbd_daemon]
command=/opt/bbd/bbd_daemon
redirect_stderr=true
autorestart=true
