# bbd_server/docker_assets/supervisor_confs/bbd_daemon.conf
[program:ttracker_server]
command=/bin/bash -c "/opt/bbd/bbd_daemon start"
redirect_stderr=true
